# crontab
SHELL=/bin/sh
20 2 * * * cd /app && python manage.py fetch_weather >> /var/log/cron.log 2>&1